# docker/base/config/tomcat-supervisor.conf
[program:tomcat]
command=/home/crafter/craftercms/apache-tomcat/bin/catalina.sh run
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
user=crafter
