# docker/authoring/config/deployer-supervisor.conf
[program:deployer]
command=/home/crafter/craftercms/deployer/deployerStart.sh
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
user=crafter
